# cmd/server/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=debug
pidfile=/tmp/supervisord.pid
nodaemon=true
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:postgres]
command=/app/init_postgres.sh
user=postgres
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/postgres.err.log
stdout_logfile=/var/log/supervisor/postgres.out.log

[program:go-server]
command=/app/server
directory=/app
user=postgres
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/go-server.err.log
stdout_logfile=/var/log/supervisor/go-server.out.log

[group:appgroup]
programs=postgres,go-server
priority=999